# overlay/etc/s6-overlay/s6-rc.d/mariadb-log/finish
#!/usr/bin/nu

use bf
use bf-s6
bf env load -x mariadb-log

# Stop log forward service
def main [...args] { bf-s6 svc finish }
