# overlay/etc/s6-overlay/s6-rc.d/mariadb-log/run
#!/command/with-contenv bash

set -euo pipefail
export BF_E="${PWD##*/}/$(basename ${0})"


#======================================================================================================================
# Forward MariaDB log to Docker.
#======================================================================================================================

bf-forward mariadb-log ${MARIADB_ERROR_LOG}
