# overlay/etc/s6-overlay/s6-rc.d/mariadb-log/run
#!/usr/bin/nu

use bf
use bf-s6
bf env load -x mariadb-log

# Forward MariaDB log to Docker
def main [...args] { bf-s6 log forward mariadb-log (bf env DB_ERROR_LOG) }
