# overlay/etc/periodic/monthly/force-upgrade
#!/bin/withenv sh

set -euo pipefail


#======================================================================================================================
# Remove version files.
#======================================================================================================================

rm -rf ${BF_DB_DATA_VERSION_FILE_OLD}
rm -rf ${BF_DB_DATA_VERSION_FILE_NEW}


#======================================================================================================================
# Run upgrade executable.
#======================================================================================================================

db-upgrade
