# overlay/etc/periodic/monthly/force-upgrade
#!/bin/withenv sh

set -euo pipefail


#======================================================================================================================
# Remove upgrade info files so the upgrade runs.
#======================================================================================================================

bf-write-debug "Removing version files."
rm -rf `bf-env DB_DATA_VERSION_FILE_OLD`
rm -rf `bf-env DB_DATA_VERSION_FILE_NEW`


#======================================================================================================================
# Run db-upgrade executable.
#======================================================================================================================

db-upgrade
