# overlay/etc/periodic/monthly/force-upgrade
#!/usr/bin/nu

use bf
use bf-mariadb
bf env load -x

# Remove upgrade files and execute
export def main [] {
    # once version files are removed, upgrade script will always run
    bf write debug "Removing version files."
    rm -rf (bf env DB_DATA_VERSION_FILE_OLD)
    rm -rf (bf env DB_DATA_VERSION_FILE_NEW)

    # execute upgrade
    bf-mariadb upgrade
}
