# overlay/etc/periodic/monthly/force-upgrade
#!/command/with-contenv bash

set -euo pipefail
export BF_E=`basename ${0}`


#======================================================================================================================
# Remove upgrade info file so the upgrade runs.
#======================================================================================================================

bf-debug "Removing ${BF_DB_UPGRADE_INFO}."
rm -f ${BF_DB_UPGRADE_INFO}


#======================================================================================================================
# Run db-upgrade executable.
#======================================================================================================================

db-upgrade
