# overlay/etc/periodic/monthly/force-upgrade
#!/bin/withenv sh

set -euo pipefail


#======================================================================================================================
# Run upgrade executable.
#======================================================================================================================

db-upgrade
