# overlay/etc/periodic/monthly/force-upgrade
#!/usr/bin/with-contenv bash

set -euo pipefail
export BCG_E=`basename ${0}`


#======================================================================================================================
# Remove upgrade info file so the upgrade runs.
#======================================================================================================================

bf-debug "Removing ${MARIADB_UPGRADE_INFO}."
rm -f ${MARIADB_UPGRADE_INFO}


#======================================================================================================================
# Run db-upgrade executable.
#======================================================================================================================

db-upgrade
