# overlay/etc/s6-overlay/s6-rc.d/mariadb/finish
#!/usr/bin/nu

use bf
use bf-mariadb
use bf-s6
bf env load -x mariadb

# Stop server and terminate container
def main [...args] {
    # old PIDs clog up the mysql directory
    bf write "Clearing temporary files and old process IDs."
    rm --force $"(bf env DB_DATA)/*.err"
    rm --force $"(bf env DB_DATA)/*.pid"

    # terminate container
    bf-s6 svc finish --terminate
}
