# overlay/etc/s6-overlay/s6-rc.d/mariadb/run
#!/bin/withenv sh

set -euo pipefail


#======================================================================================================================
# Run preflight checks
#======================================================================================================================

source /usr/lib/bf/preflight.sh mariadb


#======================================================================================================================
# If an init file exists, use it.
#======================================================================================================================

if [ -f "${BF_DB_INIT_FILE}" ] ; then

    # start server with temporary init file
    bf-write "Starting mariadbd-safe with initialisation file."
    mariadbd-safe --user=dbadm --init-file=${BF_DB_INIT_FILE}

else

    # start database
    bf-write "Starting mariadbd-safe."
    mariadbd-safe --user=dbadm

fi
