# overlay/etc/services.d/mariadb-log/run
#!/usr/bin/with-contenv bash

set -euo pipefail


#======================================================================================================================
# Forward MariaDB log to Docker
#======================================================================================================================

_forward "mariadb-log" /var/log/mariadb/error.log
