# overlay/etc/s6-overlay/s6-rc.d/mariadb-upgrade/run
#!/bin/withenv sh

set -euo pipefail


#======================================================================================================================
# Run preflight checks
#======================================================================================================================

nu -c "use bf-mariadb-upgrade run ; run preflight"


#======================================================================================================================
# Run upgrade executable.
#======================================================================================================================

db-upgrade
