# overlay/etc/services.d/mariadb-upgrade/finish
#!/bin/bash

s6-svc -d /var/run/s6/services/mariadb-upgrade
