# overlay/etc/services.d/mariadb-upgrade/run
#!/bin/bash

#======================================================================================================================
# Wait until MariaDB is running before attempting the upgrade
#======================================================================================================================

if [ "$(healthcheck)" = "1" ] ; then 
    # wait 2s before exiting the service - S6 will keep restarting it until MariaDB comes online
    echo " - waiting for MariaDB to come online..."
    sleep 2s
else
    # run upgrade executable
    echo " - MariaDB online, attempting upgrade..."
    /usr/local/bin/db-upgrade

    # de-register the upgrade service
    s6-svc -d /var/run/s6/services/mariadb-upgrade
fi
