# overlay/etc/services.d/mariadb-upgrade/run
#!/bin/bash

MARIADBD=$(pidof mariadbd)
if [ -z "${MARIADBD}" ] ; then 
    echo " - waiting for MariaDB to come online..."
    sleep 2s
else
    WAIT=10s
    echo " - MariaDB online, waiting ${WAIT} before attempting upgrade..."
    sleep ${WAIT}
    /usr/local/bin/db-upgrade
fi
