# overlay/etc/services.d/mariadb-upgrade/run
#!/usr/bin/with-contenv bash

set -euo pipefail


#======================================================================================================================
# Wait until MariaDB is running before attempting the upgrade
#======================================================================================================================

if [ -d "/var/lib/mysql/mysql" ] && [ ! -z "$(pidof mariadbd)" ]; then

    # run upgrade executable
    bcg-echo "$(mariadb --version)"
    db-upgrade

    # disable the upgrade service
    db-upgrade-disable

else

    # wait 2s before exiting the service - S6 will keep restarting it until MariaDB comes online
    # on first run, the mariadb service will disable this upgrade service itself
    SLEEP=2
    bcg-echo "Waiting ${SLEEP}s for MariaDB to come online..."
    sleep ${SLEEP}

fi
