# overlay/etc/s6-overlay/s6-rc.d/init/up
foreground
{
    ifelse -X { bf-init }
    {
        bf-write-ok "Initialisation complete." init/up
    }

    bf-cont-terminate init/up
}
