# overlay/etc/services.d/mariadb/run
#!/usr/bin/with-contenv bash

set -euo pipefail


#======================================================================================================================
# If the mysql directory exists, the server has already been initialised
#======================================================================================================================

if [ -d "/var/lib/mysql/mysql" ] ; then

    # start database
    _echo "Starting mariadbd-safe"
    mariadbd-safe


#======================================================================================================================
# Otherwise:
#   - check for root password
#   - check for application username and password
#   - run server installation
#   - generate initialisation file
#   - run server using initialisation file
#======================================================================================================================

else

    # disable upgrade service
    _echo "Database server not yet installed"
    db-upgrade-disable

    # check for root password
    [[ -z "${MARIADB_ROOT_PASSWORD}" ]] && _error "MARIADB_ROOT_PASSWORD must be set."

    # check for application user and password
    [[ -z "${MARIADB_USERNAME}" || -z "${MARIADB_PASSWORD}" ]] && _error "MARIADB_USERNAME and MARIADB_PASSWORD must be set."

    # run server installation
    _echo "Running mariadb-install-db..."
    mariadb-install-db --user=mysql --datadir="/var/lib/mysql"
    _done

    # generate init file
    INIT_FILE=/tmp/init.sql
    _echo "Generating initialisation file ${INIT_FILE}..."
    s6-setuidgid mysql /etc/db/generate-init-file "${INIT_FILE}"
    _done

    # start server with temporary init file
    _echo "Starting mariadbd-safe with initialisation file"
    mariadbd-safe --init-file=${INIT_FILE}

fi
