# overlay/etc/services.d/mariadb/run
#!/usr/bin/with-contenv bash

#======================================================================================================================
# Copyright (c) 2015 Thomas Boerger <https://webhippie.de>
# Originally from https://github.com/dockhippie/mariadb/blob/master/latest/overlay/etc/s6/mariadb/run
#
# Modifications copyright (c) 2021 Ben Green <https://bcgdesign.com>
#======================================================================================================================

if [ ! -d "/var/lib/mysql/mysql" ] ; then

    # disable database upgrade
    _echo "First run: disable upgrade service"
    db-upgrade-disable

    # check for MARIADB_ROOT_PASSWORD
    if [ -z "${MARIADB_ROOT_PASSWORD}" ] ; then
        _error "Database is uninitialized and MARIADB_ROOT_PASSWORD not set."
        source "/etc/services.d/mariadb/finish"
        exit 1
    fi

    # run base install
    _echo "Running mariadb-install-db..."
    mariadb-install-db --user=mysql --datadir="/var/lib/mysql"
    _done

    # create temporary init file
    INIT_FILE="/tmp/mariadb-boot.sql"

    # grant root access
    echo "" >| ${INIT_FILE}
    echo "ALTER USER 'root'@'localhost' IDENTIFIED BY '${MARIADB_ROOT_PASSWORD}';" >> ${INIT_FILE}
    echo "GRANT ALL ON *.* TO 'root'@'localhost WITH GRANT OPTION;" >> ${INIT_FILE}
    echo "DROP DATABASE IF EXISTS test;" >> ${INIT_FILE}
    echo "FLUSH PRIVILEGES;" >> ${INIT_FILE}

    # if specified, add database and user
    if [[ -n "${MARIADB_USERNAME}" && -n "${MARIADB_DATABASE}" ]] ; then
        DATABASES=$(echo ${MARIADB_DATABASE} | tr "," "\n")

        for DATABASE in ${DATABASES} ; do
            echo "CREATE DATABASE IF NOT EXISTS \`${DATABASE}\`;" >> ${INIT_FILE}
            echo "GRANT ALL PRIVILEGES ON \`${DATABASE}\`.* TO '${MARIADB_USERNAME}'@'%' IDENTIFIED BY '${MARIADB_PASSWORD}';" >> ${INIT_FILE}
        done

        echo "FLUSH PRIVILEGES;" >> ${INIT_FILE}
    fi

    # if database is not specified, use username as the default
    if [[ -n "${MARIADB_USERNAME}" && -z "${MARIADB_DATABASE}" ]] ; then
        echo "CREATE DATABASE IF NOT EXISTS \`${MARIADB_USERNAME}\`;" >> ${INIT_FILE}
        echo "GRANT ALL PRIVILEGES ON \`${MARIADB_USERNAME}\`.* TO '${MARIADB_USERNAME}'@'%' IDENTIFIED BY '${MARIADB_PASSWORD}';" >> ${INIT_FILE}
        echo "FLUSH PRIVILEGES;" >> ${INIT_FILE}
    fi

    # start database with temporary init file
    _echo "Starting mariadbd-safe with init file"
    mariadbd-safe --init-file=${INIT_FILE}

else

    # start database
    _echo "Starting mariadbd-safe"
    mariadbd-safe

fi
